# arch/systemd/bajek_redshift.timer
[Unit]
Description=Run bajek_redshift every 10min and 5sec after boot.

[Timer]
OnBootSec=5s
OnUnitActiveSec=10m

[Install]
WantedBy=timers.target
